# dev/erlcron/Old/crontab
52 22 28 1 2 greet hello raj danny
5 * * * 1 greet hello srini
* 10 * 4 * greet hello raj danny srini